# config/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:uvicorn]
command=uvicorn backend.app.main:app --host 0.0.0.0 --port 7860 --log-level info
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
startretries=3
stopsignal=INT
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:telegram_bot]
command=python -u bot/telegram_bot.py
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
startretries=3
stopsignal=INT
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:session_watcher]
command=python -u core/session_manager.py --watch
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
startretries=3
stopsignal=INT
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
