# config/supervisord.conf
[supervisord]
nodaemon=true

[program:uvicorn]
command=uvicorn backend.app.main:app --host 0.0.0.0 --port 7860 --log-level info
directory=/app
autostart=true
autorestart=true
startretries=3
stopsignal=INT

[program:telegram_bot]
command=python -u bot/telegram_bot.py
directory=/app
autostart=true
autorestart=true
startretries=3
stopsignal=INT

[program:session_watcher]
command=python -u core/session_manager.py --watch
directory=/app
autostart=true
autorestart=true
startretries=3
stopsignal=INT
[supervisord]
nodaemon=true

[program:fastapi]
command=uvicorn backend.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
stdout_logfile=/logs/fastapi.log
stderr_logfile=/logs/fastapi.err.log

[program:telegram-bot]
command=python -u bot/bot.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/logs/bot.log
stderr_logfile=/logs/bot.err.log

[program:session-manager]
command=python -u core/session_manager.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/logs/session_manager.log
stderr_logfile=/logs/session_manager.err.log
