# config/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:uvicorn]
command=uvicorn backend.app.main:app --host 0.0.0.0 --port 7860 --log-level info
directory=/app
autostart=true
autorestart=true
startretries=3
stopsignal=INT

stdout_logfile=/data/logs/uvicorn.log
stderr_logfile=/data/logs/uvicorn.err.log

[program:telegram_bot]
command=python -u bot/telegram_bot.py
directory=/app
autostart=true
autorestart=true
startretries=3
stopsignal=INT

stdout_logfile=/data/logs/telegram_bot.log
stderr_logfile=/data/logs/telegram_bot.err.log

[program:session_watcher]
command=python -u core/session_manager.py --watch
directory=/app
autostart=true
autorestart=true
startretries=3
stopsignal=INT

stdout_logfile=/data/logs/session_watcher.log
stderr_logfile=/data/logs/session_watcher.err.log
